# _develonly/tools/vm.crontab
* * * * * /usr/bin/bash -c "if ! mount | grep -q host ; then ~/bin/pa_mount_shared_folder.sh ; fi"
